# webapp-setup/crontab
*/3 * * * * python /app/webapp/maintain_tasks.py
